# deployment/netra_backend.service
[Unit]
Description=Netra High-Performance Go Backend Service
After=network.target mysql.service

[Service]
Type=simple
User=netra
Group=netra

# Define the working directory and absolute execution path
WorkingDirectory=/var/www/netra-api
ExecStart=/var/www/netra-api/netra-api

# Environment constraints
# Environment="PORT=8080" "GIN_MODE=release"

# Reliability rules
Restart=always
RestartSec=3

# Security & Limits
# Prevents backend from hoarding RAM in case of unexpected memory leaks
MemoryMax=512M
# Required for high concurrency WebSocket active connections
LimitNOFILE=65536

# Standard Output Redirection
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netra-backend

[Install]
WantedBy=multi-user.target
